# user/inittab
gtt1:2:respawn:/sbin/getty /dev/tty0
this-is-intentionally:malformed--dsa;;;@$$
